# roles/ndppd/files/ndppd.service
[Unit]
Description=NDP Proxy daemon

[Service]
ExecStart=/usr/sbin/ndppd

[Install]
WantedBy=default.target
